# deploy/valo-tui-ssh.service
# The public SSH front door (Wish). Reads the same SQLite cache the fetcher
# writes — read-only, so many concurrent users never contend with the writer.
#
#   sudo cp deploy/valo-tui-ssh.service /etc/systemd/system/
#   sudo systemctl enable --now valo-tui-ssh
#
# This unit binds :22 so users connect with a bare `ssh valo.example.com`. That
# requires moving the host's own admin sshd to another port FIRST (e.g. 2222) —
# see deploy/README.md — otherwise this service can't claim :22. To instead run
# on an unprivileged port, set Environment=VALO_TUI_SSH_PORT=23234 and drop the
# AmbientCapabilities line below (users then connect with `ssh -p 23234 …`).
[Unit]
Description=valo-tui SSH server (Wish)
After=network-online.target
Wants=network-online.target

[Service]
User=valo
WorkingDirectory=/var/lib/valo-tui
Environment=VALO_TUI_DB=/var/lib/valo-tui/cache.db
Environment=VALO_TUI_SSH_PORT=22
Environment=VALO_TUI_SSH_HOST_KEY=/var/lib/valo-tui/.ssh/id_ed25519
ExecStart=/usr/local/bin/valo-tui-ssh
Restart=always
RestartSec=10

# Bind the privileged port 22 as the unprivileged `valo` user. Ambient caps are
# granted even with NoNewPrivileges, and the bounding set restricts the process
# to exactly this one capability — nothing else.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening — read-only TUI; it needs its host key + the cache.
StateDirectory=valo-tui
ProtectSystem=strict
ReadWritePaths=/var/lib/valo-tui
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
